# systemd/road-photo-capture.service
[Unit]
Description=Road Photo Capture Service
Documentation=https://github.com/your-repo/road-photo-capture
After=network.target

[Service]
Type=simple
User=roadcapture
Group=roadcapture
WorkingDirectory=/opt/road-photo-capture
ExecStart=/opt/road-photo-capture/venv/bin/python /opt/road-photo-capture/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/etc/road-photo-capture/.env

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
